# node-whois-server.conf.sample
[program:node-whois-server]
command=/usr/bin/node /home/thaddeus/node-whois-server/whois-server.js
directory=/home/thaddeus/node-whois-server
numprocs=1
